# services/crontab
* * * * * (cd /app && php artisan schedule:run -q)
